# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:trading_bot]
command=python alphaalgo_2_0.py
directory=/app
user=trading
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/trading_bot.log
stderr_logfile=/var/log/trading_bot_error.log
environment=
    PYTHONPATH="/app",
    CUDA_VISIBLE_DEVICES="0"

[program:api]
command=uvicorn api:app --host 0.0.0.0 --port 8000 --workers 4
directory=/app
user=trading
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/api.log
stderr_logfile=/var/log/api_error.log

[program:monitoring]
command=python -m infrastructure.monitoring
directory=/app
user=trading
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/monitoring.log
stderr_logfile=/var/log/monitoring_error.log

[program:health_check]
command=python -m infrastructure.health_check
directory=/app
user=trading
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/health_check.log
stderr_logfile=/var/log/health_check_error.log

[program:auto_scaling]
command=python -m infrastructure.auto_scaling
directory=/app
user=trading
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/auto_scaling.log
stderr_logfile=/var/log/auto_scaling_error.log

[program:prometheus]
command=prometheus --config.file=/etc/prometheus/prometheus.yml
user=prometheus
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/prometheus.log
stderr_logfile=/var/log/prometheus_error.log

[program:grafana]
command=grafana-server --config=/etc/grafana/grafana.ini
user=grafana
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/grafana.log
stderr_logfile=/var/log/grafana_error.log

[group:alphaalgo]
programs=trading_bot,api,monitoring,health_check,auto_scaling,prometheus,grafana
priority=999

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
